# recon/cronjobs.cron
#Scan periodicaly all ports from the IP ranges in target.txt
1 */4 * * * ([ -f /var/run/nmap.running] || touch /var/run/nmap.running; nmap -P0 -p- -oX /root/current.xml -oG /root/current.txt -iL target.txt >/dev/null; rm /var/run/nmap.running)
# Perform diff on scan result periodicaly and post the result in slack 
15 */4 * * * (while [ -f /var/run/nmap.running ]; do sleep 5; done; diff /root/baseline.txt /root/current.txt | xargs -0 ./postSlack.sh ; cp /root/baseline.txt "/root/history_baseline/baseline_$(date +\%F_\%R).txt";mv /root/current.txt /root/baseline.txt)
# Periodicaly perform a DNS squatting check with dnstwist and post the result in Slack
40 * * * * (python /opt/dnstwist/dnstwist.py --registered  google.com --csv --dictionary /opt/dnstwist/dictionaries/english.dict | xargs -0 ~/postSlack.sh)
